# backend/legacy_scalping/scripts/com.ai-trade.stop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-trade.stop</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/beye82/Workspace/ai-trade/scripts/stop_simulation.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/beye82/Workspace/ai-trade</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- 월~금 15:00 종료 -->
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
            <key>Weekday</key><integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
            <key>Weekday</key><integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
            <key>Weekday</key><integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
            <key>Weekday</key><integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
            <key>Weekday</key><integer>5</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/beye82/Workspace/ai-trade/logs/launchd-stop.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/beye82/Workspace/ai-trade/logs/launchd-stop.log</string>
</dict>
</plist>
